# src/liveinternet/certbot/crontab
# Запуск обновления сертификатов в 2:30 утра каждый день
30 1 * * * /usr/local/bin/renew_certificates.sh >> /var/log/cron.log 2>&1
