# src/liveinternet/parser/crontab
55 20 * * * root . /etc/environment; /usr/local/bin/python3 /app/parser.py >> /var/log/cron.log 2>&1
