# src/liveinternet/parser/crontab
# Это файл cron, который запускает скрипт parser.py каждый день в 23:55 по МСК
55 20 * * * /usr/local/bin/python3 /app/parser.py >> /var/log/cron.log 2>&1

