# src/liveinternet/parser/crontab
# Это задание cron, которое запускает скрипт wrapper.sh каждый день в 23:55 по МСК
55 20 * * * /app/wrapper.sh
