# src/liveinternet/parser/crontab
55 20 * * * root . /etc/environment; python /app/parser.py >> /var/log/cron.log 2>&1
