# src/liveinternet/parser/crontab
55 20 * * * /bin/bash -c 'source /app/env.sh && /usr/bin/python3 /app/parser.py'
